# deploy/systemd/purl.service
[Unit]
Description=Purl Log Aggregation Server
Documentation=https://github.com/your-username/purl
After=network.target clickhouse-server.service
Wants=clickhouse-server.service

[Service]
Type=simple
User=purl
Group=purl
WorkingDirectory=/opt/purl

# Environment file
EnvironmentFile=/etc/purl/purl.env

# Perl library path
Environment=PERL5LIB=/opt/purl/lib

# Start command
ExecStart=/usr/bin/perl /opt/purl/bin/purl-server

# Restart on failure
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=purl

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/purl/data

[Install]
WantedBy=multi-user.target
