# deploy/systemd/vector.service
[Unit]
Description=Vector Log Collector for Purl
Documentation=https://vector.dev
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
Group=root

# Environment file
EnvironmentFile=-/etc/default/vector

# Start command
ExecStart=/usr/bin/vector --config /etc/vector/vector.toml
ExecReload=/bin/kill -HUP $MAINPID

# Restart on failure
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vector

# Allow access to Docker socket
SupplementaryGroups=docker

[Install]
WantedBy=multi-user.target
